# files/etc/hotplug.d/net/30-wifi-up
#!/bin/sh

if [ "$ACTION" = "add" ] && [ "$INTERFACE" = "wlan0" ] || [ "$INTERFACE" = "wlan1" ]; then
        /etc/Loglevel-helper-for-hostapd
fi

